# templates/getty@tty1.service.j2
[Service]
ExecStart=
ExecStart=-/sbin/agetty -a {{ retropie_user }} --noclear %I $TERM
